# cups/rootfs/etc/services.d/cups/run
#!/command/with-contenv bashio
# rootfs/etc/services.d/cups/run
# Start CUPS daemon after dependencies are ready

set -e

bashio::log.info "Starting CUPS daemon..."

# Double-check persistent config symlink (safety if cont-init missed it)
if [ ! -L /etc/cups ] && [ -d /data/cups/config ]; then
    rm -rf /etc/cups || true
    ln -s /data/cups/config /etc/cups || true
fi

# Start CUPS in foreground with s6 notification on port availability
exec s6-notifyoncheck -d -n 300 -w 1000 \
    -c "nc -z localhost 631" \
    /usr/sbin/cupsd -f
