# cups/rootfs/etc/services.d/airprint-watcher/run
#!/command/with-contenv bashio
# Watch for CUPS printer changes and regenerate Avahi AirPrint services

set -e

AVAHI_SERVICES_DIR="/etc/avahi/services"
CUPS_CONFIG_DIR="/data/cups/config"

bashio::log.info "Starting AirPrint service watcher..."

# Ensure services directory exists
mkdir -p "${AVAHI_SERVICES_DIR}"

# Wait for CUPS to be fully ready
sleep 5

# Initial generation of AirPrint services
bashio::log.info "Generating initial AirPrint service files..."
/usr/local/bin/airprint-generate.py 2>&1 || bashio::log.warning "Initial AirPrint generation failed (no printers configured yet?)"

# Watch for printer config changes
bashio::log.info "Watching ${CUPS_CONFIG_DIR} for printer changes..."
exec /usr/bin/inotifywait -m -e close_write,moved_to,create "${CUPS_CONFIG_DIR}" 2>&1 | \
while read -r directory events filename; do
    if [ "$filename" = "printers.conf" ]; then
        bashio::log.info "Printer config changed (${events}), regenerating AirPrint services..."
        # Small delay to ensure CUPS has finished writing
        sleep 2
        # Remove old AirPrint service files
        rm -f "${AVAHI_SERVICES_DIR}"/AirPrint-*.service
        # Regenerate service files
        /usr/local/bin/airprint-generate.py 2>&1 || bashio::log.warning "AirPrint generation failed"
        bashio::log.info "AirPrint services regenerated"
    fi
done
