# cups/rootfs/etc/services.d/avahi/run
#!/bin/sh
# rootfs/etc/services.d/avahi/run
# Run/monitor Avahi. Some avahi builds do not offer a no-fork foreground mode,
# so this script launches Avahi and monitors the daemon process, restarting if needed.

set -e

# loop forever; start avahi and then wait while it's present
while true; do
  # Try to run avahi in foreground if supported, otherwise start daemonized
  if avahi-daemon --help 2>&1 | grep -q -- '--no-.*daemon\|--no-daemon'; then
    # attempt foreground/no-daemon mode
    avahi-daemon --no-daemon || true
  else
    # fallback: daemonize and then monitor process
    avahi-daemon -D || true
  fi

  # Wait until process dies; when it dies loop restarts it
  while pgrep -x avahi-daemon > /dev/null 2>&1; do
    sleep 2
  done

  echo "avahi-daemon exited; restarting in 2s..." >&2
  sleep 2
done
